# matcha-bot.service
[Unit]
Description=Need Mo Matcha Bot
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/NeedMoMatchaBot
ExecStart=/usr/bin/python3 /home/pi/NeedMoMatchaBot/bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 